# linux-services/logistics-admin.service
[Unit]
Description=Logistics Admin Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/logistics-admin
ExecStart=/var/www/logistics-admin/Logistics.AdminApp
Restart=always
RestartSec=30
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://*:7001

[Install]
WantedBy=multi-user.target